# site/maptaxonomy/supervisor.conf
[program:uwsgi]
command=/usr/local/bin/uwsgi --socket /tmp/uwsgi.sock --chdir /srv/www/maptaxonomy/site/maptaxonomy/ --module maptaxonomy.wsgi --chmod-socket=666
autostart=true
autorestart=true
stdout_logfile=/var/uwsgi.log
redirect_stderr=true
stopsignal=QUIT
[supervisord]
logfile = /tmp/supervisord.log
logfile_maxbytes = 50MB
logfile_backups=10
loglevel = info
pidfile = /tmp/supervisord.pid
nodaemon = false
minfds = 1024
minprocs = 200
umask = 022
identifier = supervisor
directory = /tmp
nocleanup = true
childlogdir = /tmp
strip_ansi = false
environment = KEY1="value1",KEY2="value2"